# system_files/etc/systemd/system/houston-dbus.service
[Unit]
Description=Start houston dbus server

[Service]
ExecStart=/opt/45drives/houston/dbus_server.py

[Install]
WantedBy=multi-user.target
